# data/gnome-initial-setup-copy-worker.service.in
[Unit]
Description=GNOME Initial Setup Copy Worker

# Make sure we run really early
Before=gnome-session-pre.target graphical-session-pre.target

# Never run in GDM
ConditionUser=!@system

ConditionPathExists=!%E/gnome-initial-setup-done

[Service]
Type=oneshot
ExecStart=@libexecdir@/gnome-initial-setup-copy-worker
Restart=no
