# data/gnome-initial-setup-copy-worker.service.in
[Unit]
Description=GNOME Initial Setup Copy Worker

# Make sure we complete very early before most consumers are started.
Before=default.target graphical-session-pre.target

# Run before any systemd activated consumers.
Before=gnome-keyring-daemon.service

# Never run in GDM
ConditionUser=!@system

ConditionPathExists=!%E/gnome-initial-setup-done

[Service]
Type=oneshot
ExecStart=@libexecdir@/gnome-initial-setup-copy-worker
Restart=no
