# data/gnome-initial-setup-first-login.service.in
[Unit]
Description=GNOME Initial Setup

BindsTo=gnome-session.target
After=gnome-session.target

Conflicts=gnome-session@gnome-initial-setup.target
ConditionPathExists=!%E/gnome-initial-setup-done

[Service]
Type=simple
ExecStart=@libexecdir@/gnome-initial-setup --existing-user
Restart=never

[Install]
WantedBy=gnome-session.target
