# data/gnome-welcome-tour.service.in
[Unit]
Description=GNOME Welcome Tour

BindsTo=gnome-session.target
After=gnome-session.target

ConditionPathExists=%E/run-welcome-tour

[Service]
Type=simple
ExecStart=@libexecdir@/gnome-welcome-tour
Restart=no

[Install]
WantedBy=gnome-session.target
