# data/gnome-initial-setup-copy-worker.service.in
[Unit]
Description=GNOME Initial Setup Copy Worker
DefaultDependencies=no

# Make sure we run before anything else can
Before=basic.target

ConditionPathExists=!%E/gnome-initial-setup-done

[Service]
Type=oneshot
ExecStart=@libexecdir@/gnome-initial-setup-copy-worker
Restart=no
