# data/gnome-initial-setup-copy-worker.service.in
[Unit]
Description=GNOME Initial Setup Copy Worker

# Make sure we complete very early before any socket or D-Bus activated
# consumers are started.
DefaultDependencies=no
Conflicts=shutdown.target
Before=basic.target shutdown.target

# But start after tmpfiles to honor any setup there.
After=systemd-tmpfiles-setup.service

# Never run in GDM
ConditionUser=!@system

ConditionPathExists=!%E/gnome-initial-setup-done

[Service]
Type=oneshot
ExecStart=@libexecdir@/gnome-initial-setup-copy-worker
Restart=no
# We must not try to use the D-Bus session bus before `basic.target`,
# so force use of the 'local' GVfs implementation that is safe to use.
Environment=GIO_USE_VFS=local
