# data/gnome-initial-setup-first-login.service.in
[Unit]
Description=GNOME Initial Setup

BindsTo=gnome-session.target
After=gnome-session.target

ConditionPathExists=!%E/gnome-initial-setup-done

# Never run in GDM
ConditionUser=!@system
ConditionEnvironment=XDG_SESSION_CLASS=user

[Service]
Type=oneshot
ExecStart=@libexecdir@/gnome-initial-setup --existing-user
Restart=no
